# packaging/email.service

[Unit]
Description=Start the Email service

[Service]
EnvironmentFile=/run/tizen-mobile-ui
ExecStart=/usr/bin/email-service
Nice=5
OOMScoreAdjust=100
Restart=always

[Install]
WantedBy=tizen-middleware.target
